# deploy/systemd/fitpro-backup.timer
[Unit]
Description=FitPro backup — ежедневно ночью

[Timer]
OnCalendar=*-*-* 02:47:00 UTC
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
