# deploy/systemd/fitpro-backup.service
[Unit]
Description=FitPro backup (pg_dump + uploads + secrets)

[Service]
Type=oneshot
ExecStart=/opt/fitpro/deploy-scripts/backup.sh
